# ws-websocket.service
[Unit]
Description=SSH Websocket Service AJI-STORE
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/ws-openssh 8880
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=multi-user.target

